# templates/flask.socket.j2
[Unit]
Description={{ flask_service.description }} socket

[Socket]
ListenStream={{ flask_service.socket }}
SocketUser={{ flask_service.user }}

[Install]
WantedBy=sockets.target