# templates/learners.socket.j2
[Unit]
Description=learners socket

[Socket]
ListenStream=/run/learners.sock
SocketUser={{ learners_service.user }}

[Install]
WantedBy=sockets.target